# .debinate/root/etc/supervisor/conf.d/hello.conf
[program:hello]
command=/opt/hello/bin/hello
directory=/opt/hello/
user=nobody
autostart=true
autorestart=true
stdout_logfile=/var/log/hello.log
redirect_stderr=true